# deploy/systemd/todolist.service
[Unit]
Description=Todo List Spring Boot Service
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
User=instance-user
Group=instance-group
WorkingDirectory=/opt/todolist
EnvironmentFile=-/opt/todolist/.env
ExecStart=/usr/bin/java $JAVA_OPTS -jar /opt/todolist/app.jar --spring.profiles.active=${SPRING_PROFILES_ACTIVE}
Restart=on-failure
RestartSec=5
SuccessExitStatus=143
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target